# systemd/toranoi-monitor/toranoi-monitor.timer
[Unit]
Description=Run TORANOI Service Monitor every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
