# systemd/toranoi-monitor/toranoi-monitor.service
[Unit]
Description=TORANOI Service Monitor

[Service]
Type=oneshot
ExecStart=/home/toramutton/toranoi-monitor/check_services.sh
User=toramutton
